# sarg.crontab
SHELL=/bin/sh
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
NICE=15

# Run sarg after logrotate
13 5 * * *      stats    /usr/bin/sarg > /dev/null
