# sarg.crontab
SHELL=/bin/sh
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
NICE=15

# Run sarg after logrotate
03 1 * * *      stats    /usr/bin/sarg

